# scripts/epub-watcher.service
[Unit]
Description=EPUB Optimizer Output Watcher
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/bin/epub-watcher.sh
Restart=on-failure
RestartSec=5s

# Keep the service alive even if the watch dir temporarily disappears
StandardOutput=journal
StandardError=journal
SyslogIdentifier=epub-watcher

[Install]
WantedBy=default.target