# scripts/epub-optimizer.service
[Unit]
Description=EPUB Bookdrop Optimizer
After=network.target
# Start after the output watcher so the pipeline is fully ready
After=epub-watcher.service

[Service]
Type=simple
ExecStart=%h/.local/bin/epub-optimizer.sh
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=epub-optimizer

[Install]
WantedBy=default.target